# templates/varnishncsa.service.j2
[Service]
ExecStart=
ExecStart=/usr/bin/varnishncsa -a -w {{ varnish_ncsa_logfile }} -D -c -f {{ varnish_ncsa_log_formatfile }}