# templates/varnishncsa.service.j2
[Service]
ExecStart=
ExecStart=/usr/bin/varnishncsa -a -w /var/log/varnish/varnishncsa.log -D -c -f /var/log/varnish/varnishncsa_log_format